# unlock_watcher.service
[Unit]
Description=Log Watcher for fail unlock attempts

[Service]
ExecStart=/usr/local/bin/unlock_watcher

[Install]
WantedBy=multi-user.target
